# supervisor/cwpsrvd.conf
[program:cwpsrvd]
command=/usr/local/cwpsrv/bin/cwpsrvd -D FOREGROUND
autostart=true
autorestart=false