# supervisor/httpd.conf
[program:httpd]
command=/usr/local/apache/bin/httpd -D FOREGROUND
autostart=true
autorestart=unexpected