# supervisor/vsftpd.conf
[program:vsftpd]
command=/usr//bin/vsftpd
autostart=true
autorestart=unexpected