# supervisor/vsftpd.conf
[program:vsftpd]
command=/usr/sbin/vsftpd
autostart=true
autorestart=unexpected